# supervisor/go-space-chat.ini
[program:go-space-chat]
directory=/www/go-space-chat/
command= go run main.go

startsecs=1                   ; 동일합니다. 프로세스가 시작된 후 성공적인 시작으로 인식되기 전에 몇 초 동안 실행됩니다. 기본값은 1입니다.
startretries=5                ; 최대 실패 시도 횟수, 기본값 3

redirect_stderr=true
stdout_logfile=AUTO
